# deploy/bishopric.socket
[Unit]
Description=%(short_desc)s gunicorn daemon

[Socket]
ListenStream=/run/%(short_desc)s.sock

[Install]
WantedBy=sockets.target
